# site_card/uburu_uwsgi.ini
server {
    listen 80;
    server_name ;

    location /static/ {
        root /home/uburu_site/site_card;
    }

    location / {
        uwsgi_pass unix:///home/uburu_site/site_card_uwsgi.sock;
        include uwsgi_params;
    }
}

server {
    listen 80;
    server_name 194.87.99.112;

    location /static/ {
        root /home/uburu_site/site_card;
    }

    location / {
        uwsgi_pass unix:///home/uburu/uburu_site/site_card/uwsgi/uburu_uwsgi.sock;
        include uwsgi_params;
    }
}


[program:uburu_uwsgi]
environment=PATH="/home/uburu/uburuenv/bin"
numprocs=1
directory=/home/uburu/uburu_site/site_card
command=/home/uburu/uburuenv/bin/uwsgi uwsgi/production.ini
user=uburu
autostart=true
autorestart=true
redirect_stderr=true
stopwaitsecs=60
stopsignal=INT
stderr_logfile=/home/uburu/logs/%(program_name)s_err.log
stdout_logfile=/home/uburu/logs/%(program_name)s_out.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=30
stdout_capture_maxbytes=1MB




upstream django {

        server unix:///home/uburu/run/uburu.sock
}

server {

    listen 80;
    server_name 194.87.99.112;

    access_log /home/uburu/logs/nginx/access.log; 
    error_log /home/uburu/logs/nginx/error.log;

    location /static/ {
        root /home/uburu/uburu_site/site_card/;
    }

    location / {
        uwsgi_pass django;
        inclide uwsgi_params;
    }
}

[uwsgi]
uid = uburu
base = /home/%(uid)
project = uburu_site/site_card
chdir = %(base)/%(project)
module = site_card.wsgi:application
home = %(base)/uburuenv
master = true 
processes = 5
socket = %(base)/run/uburu.sock
chown-socket = %(uid):www-data
chmod-socket = 660
vacuum = true

